# roles/nas/templates/nas-luks-key.path
# https://www.freedesktop.org/software/systemd/man/systemd.path.html

[Unit]
Description=Monitor the nas key file path

[Path]
## Workaround: Trigger isn't fired on boot.
#PathExists=/run/user/1000/workaround_file-not-found
PathExists={{ nas_key_file }}
Unit=media-nas.mount

# Requires systemd v250
# https://github.com/systemd/systemd/releases/tag/v250
TriggerLimitIntervalSec=60
TriggerLimitBurst=1

[Install]
WantedBy=paths.target
