# roles/nas/templates/nas-mount.service
[Unit]
Description=Unlock and mount NAS LUKS volumes
{% if nas_key_file_mount_unit %}
After={{ nas_key_file_mount_unit }}
BindsTo={{ nas_key_file_mount_unit }}
{% else %}
After=local-fs.target
ConditionPathExists={{ nas_key_file }}
{% endif %}

[Service]
Type=oneshot
RemainAfterExit=true
{% for i in range(nas_raid_devices | length) %}
ExecStart=/usr/sbin/cryptdisks_start nas{{ i }}
{% endfor %}
ExecStart=/usr/bin/systemctl start media-nas.mount

[Install]
WantedBy={{ nas_key_file_mount_unit | default('multi-user.target', true) }}
