# roles/games/templates/heroic-token-refresh.service.j2
# {{ ansible_managed }}
[Unit]
Description=Refresh Heroic (Epic/GOG) store login tokens

[Service]
Type=oneshot
ExecStart=%h/.local/bin/heroic-refresh-tokens
